# octoprint_usb_autoconnect.service
[Unit]
Description=OctoPrint USB auto connect

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/local/bin/octoprint_usb_autoconnect
